# recipes-devtools/boot-tf/files/boot-tf.service
[Unit]
Description=Boot log to tracing file

[Service]
ExecStart=@BINDIR@/boot-tf -b

[Install]
WantedBy=multi-user.target
